# gp-vpn-container/scripts/supervisord.conf
[unix_http_server]
file=/var/run/supervisor.sock

[supervisorctl]
serverurl=unix:///var/run/supervisor.sock

[rpcinterface:supervisor]
supervisor.rpcinterface_factory = supervisor.rpcinterface:make_main_rpcinterface

[supervisord]
nodaemon=true
logfile=/var/log/supervisord.log

[program:xvnc]
command=Xtigervnc :1 -geometry 1920x1080 -depth 24 -SecurityTypes None -AlwaysShared
user=vpnuser
environment=HOME="/home/vpnuser"
autorestart=true
priority=1

[program:fluxbox]
command=fluxbox
user=vpnuser
environment=DISPLAY=":1",HOME="/home/vpnuser"
autorestart=true
priority=2

[program:novnc]
command=websockify --web /usr/share/novnc/ 6080 localhost:5901
autorestart=true
priority=3

[program:chromium]
command=google-chrome --no-sandbox --disable-setuid-sandbox --disable-gpu --disable-dev-shm-usage --no-first-run --user-agent="Mozilla/5.0 (Macintosh; Intel Mac OS X 10_15_7) AppleWebKit/605.1.15 (KHTML, like Gecko) Version/17.5 Safari/605.1.15"
user=vpnuser
environment=DISPLAY=":1",HOME="/home/vpnuser"
autostart=false
autorestart=false

[program:microsocks]
command=/usr/local/bin/microsocks -i 0.0.0.0 -p 1080
autostart=true
autorestart=true
stdout_logfile=/var/log/microsocks.log
stderr_logfile=/var/log/microsocks_err.log
